# comfyui-service-manager.service
[Unit]
Description=ComfyUI Service Manager
Documentation=https://github.com/msyzzm/comfyui-service-manager
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Run as root to manage ComfyUI processes as root
User=root
Group=root
WorkingDirectory=/home/aznable/comfyui-service-manager

# Environment
Environment="PATH=/home/aznable/miniconda3/envs/comfy/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/home/aznable/miniconda3/envs/comfy/bin/python comfyui_service_manager.py server --port 9999
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5

# Security
# NoNewPrivileges=true  # Commented out - may interfere with port binding
# PrivateTmp=true       # Commented out - causes file access issues

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=comfyui-service-manager

# Resource limits
# LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
